# systemd/qcom-batt-guard.service
[Unit]
Description=Qualcomm Battery SOC Guard (ICL)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/qcom-batt-guard
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
